# docker/php/crontab
* * * * * /usr/local/bin/php -f /var/www/bitrix/bitrix/modules/main/tools/cron_events.php
